# Build/environment/lib/systemd/system/tbs-dbus-session.service
[Unit]
Description=TBS D-Bus Session Message Bus
Requires=dbus.service

[Service]
EnvironmentFile=/home/tbs/data/scripts/tbs_dbus_session.env
ExecStart=/usr/bin/dbus-daemon --session --nofork --address=${DBUS_SESSION_BUS_ADDRESS}

[Install]
WantedBy=multi-user.target
